# root/etc/init.d/mentohust
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

PROG=/usr/sbin/mentohust
LOGFILE="/tmp/mentohust.log"
PID_FILE="/tmp/mentohust.pid"
CONF_FILE="/etc/mentohust.conf"
IFACE="eth0"

get_conf_val() {
    local key="$1"
    local line=$(grep -i "^${key}=" "$CONF_FILE" 2>/dev/null | grep -v '^[[:space:]]*[#;]' | tail -n 1)
    local val=$(echo "$line" | cut -d= -f2)
    local cleaned_val=$(echo "$val" | sed 's/[ \t\r]//g')
    
    echo "$cleaned_val"
}

# 动态读取配置中的网卡名
if [ -f "$CONF_FILE" ]; then
    nic_val=$(get_conf_val "Nic")
    [ -n "$nic_val" ] && IFACE="$nic_val"
fi

boot() {
    local if_timeout=30 count=0
    
    echo "[$(date)] Boot check started for $IFACE" >> $LOGFILE
    
    # 等待网卡
    while [ ! -d "/sys/class/net/$IFACE" ] && [ $count -lt $if_timeout ]; do
        sleep 1; count=$((count + 1))
    done
    if [ ! -d "/sys/class/net/$IFACE" ]; then
        echo "[$(date)] ERROR: $IFACE not found" >> $LOGFILE
        return 1
    fi
    
    # 等待网线连接
    count=0
    while [ $count -lt 15 ]; do
        carrier=$(cat /sys/class/net/$IFACE/carrier 2>/dev/null)
        [ "$carrier" = "1" ] && break
        sleep 1; count=$((count + 1))
    done
    echo "[$(date)] Link detected" >> $LOGFILE
    
    # 等待 IP 地址
    echo "[$(date)] Waiting for IP on $IFACE..." >> $LOGFILE
    count=0
    ip_ready=0
    while [ $count -lt 20 ]; do
        if ifconfig "$IFACE" | grep -q "inet "; then
            ip_ready=1
            break
        fi
        sleep 1
        count=$((count + 1))
    done
    
    if [ $ip_ready -eq 1 ]; then
        echo "[$(date)] IP address ready after ${count}s. Starting..." >> $LOGFILE
    else
        echo "[$(date)] WARNING: No IP yet after 20s. Starting anyway (maybe Pre-Auth DHCP)." >> $LOGFILE
    fi
    
    start
}

start_service() {
    [ ! -d "/sys/class/net/$IFACE" ] && return 1
    USERNAME=$(get_conf_val "Username")
    if [ -z "$USERNAME" ]; then
        echo "[$(date)] ERROR: Username empty. Aborting start." >> $LOGFILE
        return 1
    fi

    procd_open_instance
    procd_set_param command /bin/sh -c "exec $PROG >> $LOGFILE 2>&1"
    procd_set_param respawn 15 3 60
    procd_set_param pidfile $PID_FILE
    procd_set_param netdev $IFACE
    procd_close_instance
}

stop_service() {
    echo "[$(date)] MentoHUST stopped" >> $LOGFILE
    rm -f $PID_FILE
}